# systemd/spns.target
[Unit]
Description=Bchat Push Notification Server services
